# dockerfiles/crontab
# m h  dom mon dow   command
20 * * * * root cd /srv/project/scrape/ && ./all_spiders.sh
# end of file
